# usr/lib/systemd/system/fake-proxmox-subscription.service
[Unit]
Description=Fake Proxmox Subscription - Service

[Service]
Type=oneshot
ExecStart=/usr/bin/fake-proxmox-subscription
